# crontab
59 * * * * root service tor restart && kill $(ps aux | grep '/arpit/main' |grep -v grep| awk '{print $2}') &&  /arpit/main &
